# cookbooks/megam_logging/templates/default/heka.service.erb
[Unit]
Description=Heka Provisioning Daemon (hekad)
Documentation=https://github.com/megamsys/heka

[Service]
Type=simple
ExecStart=/usr/bin/hekad -config <%= node['heka']['conf'] %>

[Install]
WantedBy=multi-user.target
