# cookbooks/megam_nodejs_server/templates/default/nodejs.conf.erb
#!upstart
description "node.js server"
author      "Megam Systems"

start on startup
stop on shutdown

script
    export HOME="/root"

    echo $$ > /var/run/tap.pid
    exec sudo /usr/local/bin/node <%= node['nodejs']['js-file'] %> >> /var/log/upstart/tap.log 2>&1
end script

pre-start script
    # Date format same as (new Date()).toISOString() for consistency
    echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Starting" >> /var/log/upstart/tap.log
end script

pre-stop script
    rm /var/run/tap.pid
    echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Stopping" >> /var/log/upstart/tap.log
end script
