# cookbooks/megam_play_server/templates/default/play-init.conf.erb
#!upstart
description "megam play server"
author      "Megam Systems"

start on startup
stop on shutdown

script
    export HOME="/root"

    echo $$ > /var/run/play.pid
    cd /usr/share/megam_play/bin/
    exec sudo ./megam_play -Dconfig.file=/usr/share/megam_play/conf/application-production.conf >> /var/log/play.sys.log 2>&1
end script

pre-start script
    # Date format same as (new Date()).toISOString() for consistency
    echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Starting" >> /var/log/play.sys.log
end script

pre-stop script
    rm /var/run/play.pid
    echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Stopping" >> /var/log/play.sys.log
end script
